# build/client/doc/bullet.service
[Unit]
Description=bullet service
Wants=systemd-networkd-wait-online.service
After=systemd-networkd-wait-online.service
[Service]
ExecStart=/opt/bullet/bin/bullet
Restart=on-failure
[Install]
WantedBy=multi-user.target